# nestor_smart_home_bridge/rootfs/etc/services.d/nestor-bridge/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start Nestor Bridge
# ==============================================================================
set -e

bashio::log.info "Starting Nestor Bridge..."
bashio::log.info "Working directory: /app"
bashio::log.info ".NET runtime: $(dotnet --version)"

cd /app
exec dotnet /app/NestorBridge.dll
